# deploy/sunray-batman.container
# Sunray Quadlet container unit for batman
#
# Deploy to mower:
#   scp sunray-batman.container batman:/etc/containers/systemd/sunray.container
#   ssh batman "systemctl daemon-reload && systemctl start sunray"
#
# Replace GITHUB_USER below with your GitHub username.
# Image is built by GitHub Actions and pushed to ghcr.io.

[Unit]
Description=Sunray mower firmware (batman)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/GITHUB_USER/sunray-batman:latest
ContainerName=sunray
Network=host
SecurityLabelDisable=true
PodmanArgs=--privileged
Volume=sunray-data:/opt/sunray/data
Environment=TZ=Europe/Copenhagen
AutoUpdate=registry

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
